# hereweare.service.in
[Unit]
Description=Presence Detection for HA.

[Service]
ExecStart=@BINDIR@/hereweare

[Install]
WantedBy=multi-user.target
